# ansible/ansible-demo/roles/nginx/templates/supervisor.conf
[program:nginx]
command={{ prefix }}/bin/nginx -p {{ prefix }} -c {{ prefix }}/etc/nginx/nginx.conf -g "daemon off;"
directory={{ prefix }}
